# systemd/git-deploy.service
[Unit]
Description=Git repo deployer

[Service]
Type=simple
ExecStart=/opt/git-deploy/git-deploy.service.sh

[Install]
WantedBy=multi-user.target
