# systemd/git-deploy.service
[Unit]
Description=Git repo deployer

[Service]
Type=simple
ExecStart=/usr/local/sbin/git-deploy

[Install]
WantedBy=multi-user.target
